# alpine/package/modbus2mqtt/files/service.initd.in
#!/sbin/openrc-run
# OpenRC init script for @PKGNAME@

description="Service for @PKGNAME@"
command="/usr/bin/@PKGNAME@"
command_args="-c /config -d /data -s /ssl"
command_user="@PKGNAME@:dialout"
command_background=true
pidfile="/run/$RC_SVCNAME.pid"
output_log="/var/log/$RC_SVCNAME.log"
error_log="/var/log/$RC_SVCNAME.err"
directory="/var/lib/@PKGNAME@"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --mode 0755 /var/log
    checkpath --directory --mode 0755 --owner @PKGNAME@:dialout /var/lib/@PKGNAME@
    checkpath --directory --mode 0755 --owner @PKGNAME@:dialout /var/lib/@PKGNAME@/.config/git
    checkpath --file-truncate --mode 0644 --owner @PKGNAME@:dialout /var/log/$RC_SVCNAME.log
    checkpath --file-truncate --mode 0644 --owner @PKGNAME@:dialout /var/log/$RC_SVCNAME.err
    checkpath --directory --mode 0755 --owner @PKGNAME@:dialout /config/@PKGNAME@
    checkpath --directory --mode 0755 --owner @PKGNAME@:dialout /data/public
    checkpath --directory --mode 0755 --owner @PKGNAME@:dialout /config/@PKGNAME@
    checkpath --file --mode 0600 --owner @PKGNAME@:dialout /ssl/secrets.txt
    # Git is used within this service. It can be confused about permissions
    git config --system --add safe.directory /data/public
}
