# alpine/repository/modbus2mqtt/files/modbus2mqtt.initd
#!/sbin/openrc-run
name="modbus2mqtt"
description="Modbus2Mqtt Web Server"
command="/usr/bin/modbus2mqtt"
command_args="${MODBUS2MQTT_OPTS:-}"
pidfile="/run/modbus2mqtt.pid"
output_log="/var/log/modbus2mqtt.log"
error_log="/var/log/modbus2mqtt.err"

depend() {
  need net
}

start_pre() {
  mkdir -p /run
}
